# system/systemd/network.service
[Unit]
Description=Network monitor

[Service]
ExecStart=/usr/local/sbin/network
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
